# launchd/ai.openclaw.hf-daily-papers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.hf-daily-papers</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/ABSOLUTE/PATH/TO/REPO/skills/hf-daily-papers/run_and_send.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>

      <!-- Required -->
      <key>TELEGRAM_TARGET</key>
      <string>YOUR_TELEGRAM_USER_ID_OR_CHANNEL</string>

      <!-- Optional -->
      <key>OPENCLAW_BIN</key>
      <string>openclaw</string>
      <key>HF_DAILY_PAPERS_PROXY</key>
      <string></string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/ai.openclaw.hf-daily-papers.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ai.openclaw.hf-daily-papers.err.log</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
